# hyprose/airootfs/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin hyprose --noclear %I $TERM
Type=simple
