# hyprose/airootfs/etc/NetworkManager/dispatcher.d/10-hyprose
#!/bin/bash
# Auto-enable WiFi on connection

interface=$1
status=$2

if [ "$status" = "up" ]; then
    # Nothing special needed, NetworkManager handles it
    :
fi
